# d3RepoMan/d3RepoManPkgRoot/Library/LaunchAgents/com.pixar.d3.RepoMan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pixar.d3.RepoMan</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Library/Application Support/d3/d3RepoMan.app/Contents/MacOS/d3RepoMan</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
</dict>
</plist>
